# crontab
0 8 * * * root cd /app && /usr/local/bin/node /app/index.js >> /var/log/cron.log 2>&1

