# crontab
0 8 * * * root /app/index.js >> /var/log/cron.log 2>&1

